# nbimage/services.d/sshd/run
#!/command/execlineb -P

exec /usr/sbin/sshd -D
