# hubimage/services.d/sshpiperd/run
#!/command/with-contenv sh

cd /srv/sshpiperd
exec /srv/sshpiperd/sshpiperd --drop-hostkeys-message /srv/sshpiperd/sshmux
