# hubimage/services.d/jupyterhub/run
#!/command/with-contenv sh

cd /srv/jupyterhub
exec jupyterhub
